# .config/systemd/user/mbsync.service
[Unit]
Description=Mailbox sync service

[Service]
Type=oneshot
ExecStart=/home/pavel/.local/bin/syncmail

[Install]
WantedBy=default.target
